# onprem-host/ansible/roles/custodial-copy/templates/files/poduser-dotfiles/config/systemd/user/podman-auto-update.service.d/override.conf
### {{ ansible_managed }}

[Service]
Environment="HTTP_PROXY={{ proxy }}"
Environment="HTTPS_PROXY={{ proxy }}"
