# onprem-host/ansible/roles/custodial-copy/templates/files/poduser-dotfiles/config/systemd/user/podman-auto-update.timer.d/override.conf
### {{ ansible_managed }}

[Timer]
# Remove original settings of daily (midnight) update
OnCalendar=
OnUnitActiveSec=5m
RandomizedDelaySec=30
# Update immediately after reboot if an update has been missed
Persistent=true
